# docker/s6-overlay-shared/s6-rc.d/fix-permissions/up
/etc/s6-overlay/scripts/fix-permissions
