# docker/s6-overlay/s6-rc.d/initial-migration/up
/etc/s6-overlay/scripts/initial-migration