# docker/s6-overlay/s6-rc.d/composer-install/up
/etc/s6-overlay/scripts/composer-install
