# docker/scheduler/s6-overlay/s6-rc.d/cron/run
#!/command/execlineb -P
with-contenv
/etc/s6-overlay/scripts/cron