# docker/s6-overlay/s6-rc.d/symlink-roadrunner/up
/etc/s6-overlay/scripts/symlink-roadrunner