# docker/crontab
* * * * * root cd /var/www/html && s6-setuidgid webuser php artisan schedule:run >/proc/1/fd/1 2>/proc/1/fd/2
