# docker/s6-overlay/s6-rc.d/landlord-migrate/up
/etc/s6-overlay/scripts/landlord-migrate