# docker/release-automation/s6-overlay/s6-rc.d/feature-flag-purge/up
/etc/s6-overlay/scripts/feature-flag-purge