# docker/release-automation/s6-overlay/s6-rc.d/sync-schedule-monitor/up
/etc/s6-overlay/scripts/sync-schedule-monitor