# docker/s6-overlay/s6-rc.d/generate-api-docs/up
/etc/s6-overlay/scripts/generate-api-docs