# docker/s6-overlay/s6-rc.d/laravel-octane/run
#!/command/execlineb -P
with-contenv
s6-notifyoncheck -d -s 3000
/etc/s6-overlay/scripts/laravel-octane